# config/supervisor.ini
[supervisord]
user=root
group=root
nodaemon = true
logfile = /dev/null
logfile_maxbytes=0

[program:gost]
command=/bin/bash -c "/usr/bin/gost -C /etc/gost.json"
process_name=gost

[program:apt-cacher-ng]
command=/bin/bash -c "/usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng"
process_name=apt-cacher-ng
